# catkin_pix/devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.bash
# generated from rosbash/env-hooks/15.rosbash.bash.em

. "/home/rishi/catkin_pix/src/ros/tools/rosbash/rosbash"
