# catkin_pix/devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.tcsh
# generated from rosbash/env-hooks/15.rosbash.tcsh.em

. "/home/rishi/catkin_pix/src/ros/tools/rosbash/rostcsh"
