# catkin_pix/devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.fish
# generated from rosbash/env-hooks/15.rosbash.fish.em

. "/home/rishi/catkin_pix/src/ros/tools/rosbash/rosfish"
