# catkin_pix/devel/.private/rosconsole/etc/catkin/profile.d/10.rosconsole.sh
# generated from rosconsole/env-hooks/10.rosconsole.sh.develspace.em

export ROSCONSOLE_CONFIG_FILE="/home/rishi/catkin_pix/src/rosconsole/config/rosconsole.config"
