# catkin_pix/devel/etc/catkin/profile.d/10.roslaunch.sh
/home/rishi/catkin_pix/devel/.private/roslaunch/etc/catkin/profile.d/10.roslaunch.sh