# catkin_pix/devel/etc/catkin/profile.d/10.rosbuild.sh
/home/rishi/catkin_pix/devel/.private/rosbuild/etc/catkin/profile.d/10.rosbuild.sh