# catkin_pix/devel/etc/catkin/profile.d/15.rosbash.bash
/home/rishi/catkin_pix/devel/.private/rosbash/etc/catkin/profile.d/15.rosbash.bash