# catkin_pix/devel/etc/catkin/profile.d/20.transform.bash
/home/rishi/catkin_pix/devel/.private/topic_tools/etc/catkin/profile.d/20.transform.bash